# config/sprunkies-demo.service
[Unit]
Description=sprunkies Simon visual demo
Documentation=file:/home/pi/sprunkies/README.md
After=lightdm.service graphical.target
Wants=lightdm.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/sprunkies
Environment=HOME=/home/pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=SDL_VIDEODRIVER=x11
Environment=SDL_AUDIODRIVER=alsa
Environment=AUDIODEV=plughw:1,0
Environment=ALSA_CARD=Headphones
Environment=SPRUNKIES_AUDIO_DEVICE=plughw:1,0
ExecStartPre=/bin/sleep 5
ExecStart=/home/pi/sprunkies/scripts/run_demo.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical.target
